# quadlets/litellm.container
# ~/.config/containers/systemd/litellm.container
# LiteLLM proxy — main-stable (v1.82.3-stable.patch.2 as of March 2026)
# Security note: https://docs.litellm.ai/blog/security-update-march-2026

[Unit]
Description=LiteLLM proxy
After=local-fs.target network-online.target postgres.service
Wants=network-online.target

[Container]
Image=ghcr.io/berriai/litellm:main-stable

Volume=%h/.config/llm-stack/litellm-config.yaml:/app/config.yaml:ro,Z

EnvironmentFile=%h/.config/llm-stack/env
# Postgres credentials sourced from ragstack.env (overrides DATABASE_URL in env)
EnvironmentFile=%h/.config/llm-stack/ragstack.env

Environment=STORE_MODEL_IN_DB=True

Network=host
# Required: Debian-based image triggers SELinux execmem denial on Fedora 43
SecurityLabelDisable=true

Exec=--config /app/config.yaml --port 4000

ContainerName=litellm

# No curl in the Debian-minimal image — use python3 (single-quoted to avoid quadlet escaping)
HealthCmd=python3 -c 'import urllib.request; urllib.request.urlopen("http://127.0.0.1:4000/health/readiness")'
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Restart=on-failure
RestartSec=15s
TimeoutStartSec=120s

[Install]
WantedBy=default.target
