# quadlets/postgres.container
# ~/.config/containers/systemd/postgres.container
# PostgreSQL for LiteLLM

[Unit]
Description=PostgreSQL for LiteLLM
After=local-fs.target

[Container]
# No semver tags available — pinned by digest (c9s stream, 2026-04-01)
Image=quay.io/sclorg/postgresql-16-c9s@sha256:83dabb08e60b1896a8b97c0a56eabcbdc8ba4e91c70765639ea90e5e1f1a0db1

# sclorg postgres data directory (not /var/lib/postgresql/data)
Volume=litellm-db-data.volume:/var/lib/pgsql/data:Z

# sclorg images use POSTGRESQL_* env vars (not POSTGRES_*)
Environment=POSTGRESQL_USER=litellm
Environment=POSTGRESQL_PASSWORD=litellm
Environment=POSTGRESQL_DATABASE=litellm

Network=host

ContainerName=postgres

[Service]
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
