# quadlets/rag-proxy.container
# ~/.config/containers/systemd/rag-proxy.container
# RAG proxy — searches Qdrant, reranks with cross-encoder, injects context
# Managed by: systemctl --user [start|stop|restart|status] rag-proxy

[Unit]
Description=RAG proxy (Qdrant search + reranker + context injection)
After=local-fs.target network-online.target qdrant.service
Wants=network-online.target
Requires=qdrant.service

[Container]
# UBI9 Python 3.11 pinned to 2026-04-01 build
Image=registry.access.redhat.com/ubi9/python-311@sha256:8fb94e142e1093b82ca64b334264a5da8e874567d64ce4b6d614f86da3e38813

Volume=%h/git/framework-ai-stack/rag-proxy:/app:ro,Z

# Writable model cache for sentence-transformers and cross-encoder models
Volume=rag-proxy-model-cache:/opt/app-root/src/.cache:Z

# NOTE: ReadOnly=true is not compatible with pip-at-runtime installs.
# To enable read-only rootfs, build a custom image with deps pre-installed.

# Run as non-root user (default in UBI Python images, UID 1001)
User=1001

Environment=MODEL_URL=http://host.containers.internal:8080
Environment=QDRANT_URL=http://host.containers.internal:6333
Environment=QDRANT_COLLECTION=documents
Environment=EMBED_MODEL=sentence-transformers/all-mpnet-base-v2
# Reduced from 20 to 8 — reranker on CPU takes ~3s/candidate,
# so 20 candidates = 60s latency. 8 candidates ≈ 24s.
Environment=RAG_TOP_K=8
Environment=RAG_PROXY_PORT=8090

# Document store — hydrates chunk text from Postgres at query time
Environment=DOCSTORE_BACKEND=postgres
Environment=DOCSTORE_URL=postgresql://litellm:litellm@host.containers.internal:5432/litellm

# Reranker configuration
# Disabled: CPU reranking adds 24-60s latency per query. Re-enable when
# GPU access is available in this container or a lighter model is used.
Environment=RERANKER_ENABLED=false
Environment=RERANKER_MODEL=BAAI/bge-reranker-v2-m3
Environment=RERANKER_DEVICE=cpu
Environment=RERANKER_TOP_N=5

# Thinking budget — modest token budget for reasoning across sources
# without unconstrained latency
Environment=THINKING_BUDGET=1024

# HuggingFace cache location (inside the writable volume)
Environment=HF_HOME=/opt/app-root/src/.cache/huggingface
Environment=TRANSFORMERS_CACHE=/opt/app-root/src/.cache/huggingface
Environment=PYTHONPATH=/app

Network=host

Exec=bash -c "pip install --quiet -r /app/requirements.txt && python3 /app/rag-proxy.py"

ContainerName=rag-proxy
PublishPort=8090:8090

# Health check — proxy should respond to any request
HealthCmd=python3 -c "import urllib.request; urllib.request.urlopen('http://127.0.0.1:8090/v1/models')"
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=120s

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300s

# Resource limits
MemoryMax=4G
CPUQuota=200%

[Install]
WantedBy=default.target
