# quadlets/rag-watcher.container
# ~/.config/containers/systemd/rag-watcher.container
# Google Drive RAG watcher — polls a Drive folder and ingests documents into Qdrant
# Managed by: systemctl --user [start|stop|restart|status] rag-watcher

[Unit]
Description=Google Drive RAG watcher
After=local-fs.target network-online.target qdrant.service
Wants=network-online.target
Requires=qdrant.service

[Container]
Image=registry.access.redhat.com/ubi10:latest

Volume=%h/git/framework-ai-stack/rag-watcher:/app:ro,Z
Volume=%h/.config/ramalama/gdrive-sa.json:/run/secrets/gdrive-sa.json:ro,Z
Volume=rag-staging:/staging:Z

EnvironmentFile=%h/.config/llm-stack/env
Environment=GOOGLE_APPLICATION_CREDENTIALS=/run/secrets/gdrive-sa.json
Environment=STAGING_DIR=/staging
Environment=QDRANT_URL=http://host.containers.internal:6333
Environment=QDRANT_COLLECTION=documents
Environment=EMBED_MODEL=sentence-transformers/all-mpnet-base-v2

Network=host

Exec=bash -c "dnf install -y -q python3.12 python3.12-pip git-core && pip3.12 install --quiet -r /app/requirements.txt && python3.12 /app/rag-watcher.py"

ContainerName=rag-watcher

[Service]
Restart=on-failure
RestartSec=30s
TimeoutStartSec=300s

[Install]
WantedBy=default.target
