# quadlets/rag-proxy.container
# ~/.config/containers/systemd/rag-proxy.container
# RAG proxy — searches Qdrant, reranks with cross-encoder, injects context
# Managed by: systemctl --user [start|stop|restart|status] rag-proxy

[Unit]
Description=RAG proxy (Qdrant search + reranker + context injection)
After=local-fs.target network-online.target qdrant.service postgres.service
Wants=network-online.target
Requires=qdrant.service postgres.service

[Container]
# Pre-built image with deps + models baked in (see rag-proxy/Containerfile)
# Build: podman build -t localhost/rag-proxy rag-proxy/
# Falls back to base UBI9 image if not built (pip install at runtime)
Image=localhost/rag-proxy:latest

Volume=%h/git/framework-ai-stack/rag-proxy:/app:ro,Z

# Writable model cache for sentence-transformers and cross-encoder models
Volume=rag-proxy-model-cache:/opt/app-root/src/.cache:Z

# Run as non-root user (default in UBI Python images, UID 1001)
User=1001

Environment=MODEL_URL=http://host.containers.internal:8080
Environment=QDRANT_URL=http://host.containers.internal:6333
Environment=QDRANT_COLLECTION=documents
Environment=EMBED_MODEL=sentence-transformers/all-mpnet-base-v2
# Force CPU — ROCm PyTorch segfaults on gfx1151 for sentence-transformers
Environment=EMBED_DEVICE=cpu
Environment=RAG_TOP_K=20
Environment=RAG_PROXY_PORT=8090

# Document store — hydrates chunk text from Postgres at query time
Environment=DOCSTORE_BACKEND=postgres
Environment=DOCSTORE_URL=postgresql://litellm:litellm@host.containers.internal:5432/litellm

# Reranker — lightweight CPU model (22M params, sub-second on CPU)
# The heavier BAAI/bge-reranker-v2-m3 (0.6B) needs GPU or takes 60s on CPU.
# ms-marco-MiniLM-L-6-v2 trades some multilingual quality for ~100x speed.
# Swap to BAAI/bge-reranker-v2-m3 when GPU reranking is stable on gfx1151.
Environment=RERANKER_ENABLED=true
Environment=RERANKER_MODEL=cross-encoder/ms-marco-MiniLM-L-6-v2
Environment=RERANKER_DEVICE=cpu
Environment=RERANKER_TOP_N=5

# Thinking budget — modest token budget for reasoning across sources
# without unconstrained latency
Environment=THINKING_BUDGET=1024

# HuggingFace cache location (inside the writable volume)
Environment=HF_HOME=/opt/app-root/src/.cache/huggingface
Environment=TRANSFORMERS_CACHE=/opt/app-root/src/.cache/huggingface
Environment=PYTHONPATH=/app

Network=host

Exec=python3 /app/rag-proxy.py

ContainerName=rag-proxy

# Health check — proxy should respond to any request
HealthCmd=python3 -c "import urllib.request; urllib.request.urlopen('http://127.0.0.1:8090/health')"
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=120s

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300s

# Resource limits
MemoryMax=5G
CPUQuota=200%

[Install]
WantedBy=default.target
