# quadlets/rag-watcher.container
# ~/.config/containers/systemd/rag-watcher.container
# Google Drive RAG watcher — polls a Drive folder and rebuilds RAG OCI images

[Unit]
Description=Google Drive RAG watcher
After=local-fs.target network-online.target
Wants=network-online.target

[Container]
Image=docker.io/library/python:3.12-slim

Volume=%h/git/framework-ai-stack/rag-watcher:/app:ro,Z
Volume=%h/.config/ramalama/gdrive-sa.json:/run/secrets/gdrive-sa.json:ro,Z
Volume=ramalama-rag-staging:/staging:Z

EnvironmentFile=%h/.config/llm-stack/env
Environment=GOOGLE_APPLICATION_CREDENTIALS=/run/secrets/gdrive-sa.json
Environment=STAGING_DIR=/staging

Network=host

Exec=bash -c "pip install --quiet -r /app/requirements.txt && python /app/rag-watcher.py"

ContainerName=rag-watcher

[Service]
Restart=on-failure
RestartSec=30s
TimeoutStartSec=180s

[Install]
WantedBy=default.target
