# quadlets/ragdeck.container
# ~/.config/containers/systemd/ragdeck.container
# ragdeck — admin UI for the rag-suite stack
# Source: https://github.com/aclater/ragdeck
# Managed by: systemctl --user [start|stop|restart|status] ragdeck

[Unit]
Description=ragdeck — RAG suite admin UI
After=network-online.target ragpipe.service ragstuffer.service ragwatch.service

[Container]
# Pin digest after first successful build: skopeo inspect docker://ghcr.io/aclater/ragdeck:main
Image=ghcr.io/aclater/ragdeck:main
ContainerName=ragdeck

PublishPort=8092:8092

# Service URLs — host.containers.internal resolves to the host from inside Podman
Environment=RAGPIPE_URL=http://host.containers.internal:8090
Environment=RAGSTUFFER_URL=http://host.containers.internal:8091
Environment=RAGWATCH_URL=http://host.containers.internal:9090
Environment=QDRANT_URL=http://host.containers.internal:6333

# Credentials from ragstack.env
EnvironmentFile=%h/.config/llm-stack/ragstack.env

# Health check
HealthCmd=python3 -c "import urllib.request; urllib.request.urlopen('http://127.0.0.1:8092/health')"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=15s

# Security
NoNewPrivileges=true
ReadOnly=true
Tmpfs=/tmp

# Networking — needs access to all rag-suite services on host
Network=host

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
