# quadlets/ragorchestrator.container
# ~/.config/containers/systemd/ragorchestrator.container
# ragorchestrator — LangGraph supervisor agent for ragpipe
# Source: https://github.com/aclater/ragorchestrator
# Managed by: systemctl --user [start|stop|restart|status] ragorchestrator

[Unit]
Description=ragorchestrator (LangGraph supervisor agent for agentic RAG)
After=local-fs.target network-online.target ragpipe.service
Wants=network-online.target
Requires=ragpipe.service

[Container]
Image=ghcr.io/aclater/ragorchestrator:main

ContainerName=ragorchestrator

# Run as non-root user (default in UBI Python images, UID 1001)
User=1001

# Service URLs — host.containers.internal resolves to the host from inside Podman
Environment=RAGPIPE_URL=http://host.containers.internal:8090
Environment=RAGWATCH_URL=http://host.containers.internal:9090
Environment=QDRANT_URL=http://host.containers.internal:6333
Environment=RAGORCHESTRATOR_PORT=8095
# LLM endpoint — must use 127.0.0.1 not hostname to avoid IPv6 resolution
# llama-vulkan listens on IPv4 only; Fedora resolves localhost/hostnames to ::1 first
Environment=MODEL_URL=http://127.0.0.1:8080
Environment=MODEL_NAME=model.file

# Web search — disabled until TAVILY_API_KEY is configured in ragstack.env
Environment=DISABLE_WEB_SEARCH=true
Environment=LANGSMITH_TRACING=false

# Credentials sourced from ragstack.env — never committed to the repo
# Provides: RAGPIPE_ADMIN_TOKEN, DOCSTORE_URL, TAVILY_API_KEY (optional)
EnvironmentFile=%h/.config/llm-stack/ragstack.env

# Networking — needs access to all rag-suite services on host
Network=host

# Health check — UBI10 minimal lacks curl, use Python
HealthCmd=python3 -c "import urllib.request; urllib.request.urlopen('http://127.0.0.1:8095/health')"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

# Security
NoNewPrivileges=true
ReadOnly=true
Tmpfs=/tmp

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=120s

[Install]
WantedBy=default.target
