# quadlets/ragorchestrator.container
# ~/.config/containers/systemd/ragorchestrator.container
# ragorchestrator — LangGraph supervisor agent for ragpipe
# Source: https://github.com/aclater/ragorchestrator
# Managed by: systemctl --user [start|stop|restart|status] ragorchestrator

[Unit]
Description=ragorchestrator (LangGraph supervisor agent for agentic RAG)
After=local-fs.target network-online.target ragpipe.service
Wants=network-online.target
Requires=ragpipe.service

[Container]
Image=ghcr.io/aclater/ragorchestrator:main

# Run as non-root user (default in UBI Python images, UID 1001)
User=1001

Environment=RAGPIPE_URL=http://host.containers.internal:8090
Environment=RAGORCHESTRATOR_PORT=8095

# Credentials sourced from ragstack.env — never committed to the repo
# Provides: RAGPIPE_ADMIN_TOKEN
EnvironmentFile=%h/.config/llm-stack/ragstack.env

Network=host

ContainerName=ragorchestrator

# Health check
HealthCmd=curl -sf http://127.0.0.1:8095/health
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=120s

[Install]
WantedBy=default.target
