# quadlets/ragpipe.container
# ~/.config/containers/systemd/ragpipe.container
# ragpipe — searches Qdrant, reranks with cross-encoder, injects context
# Source: https://github.com/aclater/ragpipe
# Managed by: systemctl --user [start|stop|restart|status] ragpipe

[Unit]
Description=ragpipe (Qdrant search + reranker + context injection)
After=local-fs.target network-online.target qdrant.service postgres.service
Wants=network-online.target
Requires=qdrant.service postgres.service

[Container]
# ROCm variant: rocm/dev-ubuntu-24.04:7.2.1 + onnxruntime-migraphx (ADR-004 exception)
# ONNX models baked in at /opt/ragpipe/models
Image=ghcr.io/aclater/ragpipe:main-rocm

# Writable cache for MIGraphX JIT-compiled kernels (persists across restarts)
Volume=ragpipe-model-cache:/home/default/.cache/ragpipe:Z

# Routes config (local file, not in repo)
Volume=%h/.config/ragpipe/routes.yaml:/etc/ragpipe/routes.yaml:ro,Z

# Run as non-root user (default in UBI Python images, UID 1001)
User=1001

Environment=MODEL_URL=http://host.containers.internal:8080
Environment=QDRANT_URL=http://host.containers.internal:6333
Environment=QDRANT_COLLECTION=documents

# Credentials sourced from ragstack.env — never committed to the repo
# Provides: RAGPIPE_ADMIN_TOKEN, DOCSTORE_URL
EnvironmentFile=%h/.config/llm-stack/ragstack.env

# Semantic router — optional, ragpipe works without RAGPIPE_ROUTES_FILE
# Create ~/.config/ragpipe/routes.yaml to enable routing (see ragpipe docs)
Environment=RAGPIPE_ROUTES_FILE=/etc/ragpipe/routes.yaml

# Override ragpipe defaults for wider retrieval window (default: 20/5)
Environment=RAG_TOP_K=40
Environment=RERANKER_TOP_N=15

# Document store — asyncpg connection pool for hydration, LRU chunk cache
Environment=DOCSTORE_BACKEND=postgres

# AMD ROCm GPU access for ONNX Runtime inference (gfx1151)
# Force MIGraphXExecutionProvider — ROCMExecutionProvider segfaults on gfx1151
# (RDNA 4 not yet supported by ROCM EP). MIGraphX uses graph-level compilation
# and handles gfx1151 via HSA_OVERRIDE_GFX_VERSION.
Environment=RAGPIPE_DEVICE=migraphx
AddDevice=/dev/kfd
AddDevice=/dev/dri
Environment=HSA_OVERRIDE_GFX_VERSION=11.5.1
# MIOpen JIT-compiles kernels for gfx1151 on first run; cache in writable volume
Environment=MIOPEN_USER_DB_PATH=/home/default/.cache/ragpipe/miopen
# SELinux cannot label /dev/kfd — required for ROCm device access
# Upstream: https://github.com/ROCm/ROCm/issues/2983
SecurityLabelDisable=true
# MIGraphX requires shared memory for GPU graph compilation
PodmanArgs=--ipc=host

Network=host

ContainerName=ragpipe

# Health check
HealthCmd=curl -sf http://127.0.0.1:8090/health
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=120s

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300s

# Resource limits — MIGraphX graph compilation peaks higher than CPU-only ONNX
MemoryMax=4G
CPUQuota=200%

[Install]
WantedBy=default.target
