# quadlets/ragpipe.container
# ~/.config/containers/systemd/ragpipe.container
# ragpipe — searches Qdrant, reranks with cross-encoder, injects context
# Source: https://github.com/aclater/ragpipe
# Managed by: systemctl --user [start|stop|restart|status] ragpipe

[Unit]
Description=ragpipe (Qdrant search + reranker + context injection)
After=local-fs.target network-online.target qdrant.service postgres.service
Wants=network-online.target
Requires=qdrant.service postgres.service

[Container]
# Pre-built image with deps + models baked in (see github.com/aclater/ragpipe)
# Build: podman build -t localhost/ragpipe ~/git/ragpipe/
Image=localhost/ragpipe:latest

Volume=%h/git/ragpipe:/app:ro,z

# Writable model cache for fastembed ONNX models
Volume=ragpipe-model-cache:/opt/app-root/src/.cache:Z

# Run as non-root user (default in UBI Python images, UID 1001)
User=1001

Environment=MODEL_URL=http://host.containers.internal:8080
Environment=QDRANT_URL=http://host.containers.internal:6333
Environment=QDRANT_COLLECTION=documents
# fastembed uses ONNX Runtime (CPU-only, no PyTorch) — no GPU segfault risk
Environment=EMBED_MODEL=BAAI/bge-base-en-v1.5
Environment=RAG_TOP_K=40
Environment=RAG_PROXY_PORT=8090

# Document store — hydrates chunk text from Postgres at query time
Environment=DOCSTORE_BACKEND=postgres
Environment=DOCSTORE_URL=postgresql://litellm:litellm@host.containers.internal:5432/litellm

# Reranker — lightweight ONNX model (22M params, sub-second on CPU)
Environment=RERANKER_ENABLED=true
Environment=RERANKER_MODEL=Xenova/ms-marco-MiniLM-L-6-v2
Environment=RERANKER_TOP_N=15

# Thinking budget — modest token budget for reasoning across sources
# without unconstrained latency
Environment=THINKING_BUDGET=1024

# HuggingFace cache location (inside the writable volume)
Environment=HF_HOME=/opt/app-root/src/.cache/huggingface
Environment=TRANSFORMERS_CACHE=/opt/app-root/src/.cache/huggingface
Environment=PYTHONPATH=/app

Network=host

Exec=python3 /app/ragpipe.py

ContainerName=ragpipe

# Health check — proxy should respond to any request
HealthCmd=python3 -c "import urllib.request; urllib.request.urlopen('http://127.0.0.1:8090/health')"
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=120s

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300s

# Resource limits
MemoryMax=5G
CPUQuota=200%

[Install]
WantedBy=default.target
