# quadlets/ragstuffer.container
# ~/.config/containers/systemd/ragstuffer.container
# ragstuffer — polls document sources and ingests into Qdrant + Postgres
# Source: https://github.com/aclater/ragstuffer
# Managed by: systemctl --user [start|stop|restart|status] ragstuffer

[Unit]
Description=ragstuffer (document ingestion for RAG)
After=local-fs.target network-online.target qdrant.service ragpipe.service
Wants=network-online.target
Requires=qdrant.service ragpipe.service

[Container]
# Locally built — tag with git SHA at build time:
#   podman build -t localhost/ragstuffer:$(git rev-parse --short HEAD) .
#   podman tag localhost/ragstuffer:$(git rev-parse --short HEAD) localhost/ragstuffer:main
Image=localhost/ragstuffer:main

Volume=%h/git/ragstuffer:/app:ro,z
Volume=%h/git/ragpipe/ragpipe/docstore.py:/app/docstore.py:ro,z
Volume=%h/.config/ramalama/gdrive-sa.json:/run/secrets/gdrive-sa.json:ro,z
Volume=rag-staging:/staging:Z

EnvironmentFile=%h/.config/llm-stack/env
# Credentials sourced from ragstack.env — never committed to the repo
# Provides: DOCSTORE_URL, RAGSTUFFER_ADMIN_TOKEN
EnvironmentFile=%h/.config/llm-stack/ragstack.env

Environment=GOOGLE_APPLICATION_CREDENTIALS=/run/secrets/gdrive-sa.json
Environment=STAGING_DIR=/staging
Environment=QDRANT_URL=http://host.containers.internal:6333
Environment=QDRANT_COLLECTION=documents
Environment=EMBED_URL=http://host.containers.internal:8090/v1/embeddings
Environment=DOCSTORE_BACKEND=postgres

# Admin API for triggering ingestion (POST /admin/ingest-now, /admin/ingest-full)
Environment=RAGSTUFFER_ADMIN_PORT=8091

Network=host

Exec=python3.12 /app/ragstuffer.py

ContainerName=ragstuffer

HealthCmd=curl -sf http://127.0.0.1:8091/health
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Restart=on-failure
RestartSec=30s
TimeoutStartSec=300s

[Install]
WantedBy=default.target
