# quadlets/ragstuffer.container
# ~/.config/containers/systemd/ragstuffer.container
# ragstuffer — polls document sources and ingests into Qdrant + Postgres
# Source: https://github.com/aclater/ragstuffer
# Managed by: systemctl --user [start|stop|restart|status] ragstuffer

[Unit]
Description=ragstuffer (document ingestion for RAG)
After=local-fs.target network-online.target qdrant.service ragpipe.service
Wants=network-online.target
Requires=qdrant.service ragpipe.service

[Container]
# Built from github.com/aclater/ragstuffer
# Build: podman build -t localhost/ragstuffer ~/git/ragstuffer/
Image=localhost/ragstuffer:latest

Volume=%h/git/ragstuffer:/app:ro,z
Volume=%h/git/ragpipe/ragpipe/docstore.py:/app/docstore.py:ro,z
Volume=%h/.config/ramalama/gdrive-sa.json:/run/secrets/gdrive-sa.json:ro,z
Volume=rag-staging:/staging:Z

EnvironmentFile=%h/.config/llm-stack/env
Environment=GOOGLE_APPLICATION_CREDENTIALS=/run/secrets/gdrive-sa.json
Environment=STAGING_DIR=/staging
Environment=QDRANT_URL=http://host.containers.internal:6333
Environment=QDRANT_COLLECTION=documents
Environment=EMBED_URL=http://host.containers.internal:8090/v1/embeddings
Environment=DOCSTORE_BACKEND=postgres
Environment=DOCSTORE_URL=postgresql://litellm:litellm@host.containers.internal:5432/litellm

# Admin API for triggering ingestion (POST /admin/ingest-now, /admin/ingest-full)
Environment=RAGSTUFFER_ADMIN_PORT=8091
# Set a real token in ~/.config/llm-stack/env — this is a placeholder
Environment=RAGSTUFFER_ADMIN_TOKEN=change-me

Network=host

Exec=python3.12 /app/ragstuffer.py

ContainerName=ragstuffer

[Service]
Restart=on-failure
RestartSec=30s
TimeoutStartSec=300s

[Install]
WantedBy=default.target
