# quadlets/ragwatch.container
# ~/.config/containers/systemd/ragwatch.container
# ragwatch — observability aggregation for rag-suite (ragpipe + ragstuffer metrics)
# Source: https://github.com/aclater/ragwatch
# Managed by: systemctl --user [start|stop|restart|status] ragwatch

[Unit]
Description=ragwatch — observability aggregation for rag-suite
After=network-online.target ragpipe.service ragstuffer.service

[Container]
# Pin digest after first successful build: skopeo inspect docker://ghcr.io/aclater/ragwatch:main
Image=ghcr.io/aclater/ragwatch:main
ContainerName=ragwatch

PublishPort=9090:9090

# Credentials from ragstack.env
EnvironmentFile=%h/.config/llm-stack/ragstack.env

# Health check
HealthCmd=python3 -c "import urllib.request; urllib.request.urlopen('http://127.0.0.1:9090/health')"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=10s

# Security
NoNewPrivileges=true
ReadOnly=true
Tmpfs=/tmp

# Networking — needs access to ragpipe and ragstuffer metrics endpoints on host
Network=host

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
