# quadlets/ramalama.container
# ~/.config/containers/systemd/ramalama.container
# Qwen3.5-35B-A3B via RamaLama ROCm container
# Managed by: systemctl --user [start|stop|restart|status] ramalama

[Unit]
Description=RamaLama Qwen3.5-35B-A3B
After=local-fs.target network-online.target
Wants=network-online.target

[Container]
# Pinned 2026-04-04 — update: skopeo inspect docker://quay.io/ramalama/rocm:latest
Image=quay.io/ramalama/rocm:0.18.0

# GPU passthrough — /dev/kfd (compute) + /dev/dri (render)
AddDevice=-/dev/kfd
AddDevice=-/dev/dri
# SELinux container_t denies access to /dev/kfd (ROCm compute interface).
# No upstream SELinux policy module exists for /dev/kfd on Fedora 43.
# Disabling label confinement is required until ROCm ships an SELinux policy
# or Fedora adds /dev/kfd to container device access rules.
# Ref: ADR-011, https://github.com/ROCm/ROCm/issues/3987
SecurityLabelDisable=true

# GFX version override for AI Max+ 395; remove once ROCm ships gfx1151 support
Environment=HSA_OVERRIDE_GFX_VERSION=11.5.1
Environment=AMD_LOG_LEVEL=0
# Limit HW queues to reduce idle GPU busy-spin on gfx1151
Environment=GPU_MAX_HW_QUEUES=1
# KV cache quantization and flash attention — CLI flags are authoritative,
# env vars kept as fallback for ramalama image entrypoint compatibility
Environment=LLAMA_ARG_FLASH_ATTN=1
Environment=LLAMA_ARG_CACHE_TYPE_K=q4_0
Environment=LLAMA_ARG_CACHE_TYPE_V=q4_0

# Model bind-mount — RamaLama stores models under ~/.local/share/ramalama
Mount=type=bind,src=MODEL_PATH_PLACEHOLDER,target=/mnt/models/model.file,ro,Z

# --flash-attn on: explicit (env var is backup, CLI is authoritative)
# --cache-type-k/v q4_0: halves KV cache vs q8_0, frees ~10GB system RAM
# --parallel 2: doubles per-slot context (32K), halves KV cache overhead
# --ctx-size 65536: 32K per slot is sufficient for RAG (typically <8K prompt)
Exec=llama-server --model /mnt/models/model.file --host 0.0.0.0 --port 8080 \
    --ctx-size 65536 --n-gpu-layers 999 --threads 16 --parallel 2 \
    --flash-attn on --cache-type-k q4_0 --cache-type-v q4_0 --jinja

ContainerName=ramalama
PublishPort=8080:8080

HealthCmd=curl -sf http://127.0.0.1:8080/health
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=120s

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=180s

[Install]
WantedBy=default.target
