# hosts/nvidia/quadlets/ramalama.container
# ~/.config/containers/systemd/ramalama.container
# Qwen3.5-9B via RamaLama CUDA on NVIDIA GPUs (≤16 GB VRAM)
# Auto-selected when nvidia-smi is detected and VRAM < 32 GB
# Managed by: systemctl --user [start|stop|restart|status] ramalama

[Unit]
Description=RamaLama Qwen3.5-9B (CUDA)
After=local-fs.target network-online.target
Wants=network-online.target

[Container]
Image=quay.io/ramalama/cuda:latest

# NVIDIA GPU passthrough via CDI (nvidia-container-toolkit + /etc/cdi/nvidia.yaml)
AddDevice=nvidia.com/gpu=all
# Upstream CUDA image triggers SELinux execmem denial on Fedora 43.
# This is a Debian binary incompatibility, not a GPU device issue (see ADR-011).
SecurityLabelDisable=true

# Model bind-mount — ramalama stores models under ~/.local/share/ramalama
Mount=type=bind,src=MODEL_PATH_PLACEHOLDER,target=/mnt/models/model.file,ro,Z

# Qwen3.5-9B at UD-Q4_K_XL (~6 GB) fits entirely in 12 GB VRAM with room
# for KV cache. 32K context with 2 parallel slots (16K each) leaves ~5 GB
# headroom for KV cache at q4_0.
# --flash-attn on: Ada Lovelace (sm_89) and newer support flash attention
Exec=llama-server --model /mnt/models/model.file --host 0.0.0.0 --port 8080 \
    --ctx-size 32768 --n-gpu-layers 999 --threads 16 --parallel 2 \
    --flash-attn on --cache-type-k q4_0 --cache-type-v q4_0 --jinja

ContainerName=ramalama
PublishPort=8080:8080

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=180s

[Install]
WantedBy=default.target
